# example.crontab
0 */6 * * * bin/app scrape --airline=OAO
1 */6 * * * bin/app scrape --airline=DXT
2 */6 * * * bin/app scrape --airline=SHU
3 */6 * * * bin/app scrape --airline=AZO
4 */6 * * * bin/app scrape --airline=AZV
5 */6 * * * bin/app scrape --airline=RSY
6 */6 * * * bin/app scrape --airline=DRU
7 */6 * * * bin/app scrape --airline=AGU
8 */6 * * * bin/app scrape --airline=BRP
9 */6 * * * bin/app scrape --airline=AFL
10 */6 * * * bin/app scrape --airline=VGV
11 */6 * * * bin/app scrape --airline=GZP
12 */6 * * * bin/app scrape --airline=IZA
13 */6 * * * bin/app scrape --airline=KAR
14 */6 * * * bin/app scrape --airline=IAE
15 */6 * * * bin/app scrape --airline=PTK
16 */6 * * * bin/app scrape --airline=SSJ
17 */6 * * * bin/app scrape --airline=TYA
18 */6 * * * bin/app scrape --airline=PBD
19 */6 * * * bin/app scrape --airline=RKA
20 */6 * * * bin/app scrape --airline=RWZ
21 */6 * * * bin/app scrape --airline=SDM
22 */6 * * * bin/app scrape --airline=RLU
23 */6 * * * bin/app scrape --airline=NWS
24 */6 * * * bin/app scrape --airline=D2
25 */6 * * * bin/app scrape --airline=SBI
26 */6 * * * bin/app scrape --airline=AUL
27 */6 * * * bin/app scrape --airline=SVR
28 */6 * * * bin/app scrape --airline=UTA
29 */6 * * * bin/app scrape --airline=TUM
30 */6 * * * bin/app scrape --airline=SYL
31 */6 * * * bin/app scrape --airline=LLM
